# backend/config/conf.d/start_beat.conf
[program:start_beat]
command=/usr/local/bin/beat_run.sh
user=webmaster
directory=/home/webmaster/neuraldating/
stdout_logfile=/var/log/dating/star_beat.log
redirect_stderr=true
